# AppBlocker/diegoibarra.AppBlocker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>diegoibarra.AppBlocker</string>
		<key>Program</key>
		<string>/Users/diegoibarra/Developer/1_myProjects/MiniProjects/AppBlocker/myAppBlocker.py</string>
		<key>WorkingDirectory</key>
		<string>/Users/diegoibarra/Developer/1_myProjects/MiniProjects/AppBlocker</string>
		<key>KeepAlive</key>
		<true/>
		<key>ThrottleInterval</key>
		<integer>10</integer>
		<key>StandardErrorPath</key>
		<string>/Users/diegoibarra/Developer/1_myProjects/MiniProjects/AppBlocker/cache/diegoibarra.AppBlocker.err</string>
		<key>StandardOutPath</key>
		<string>/Users/diegoibarra/Developer/1_myProjects/MiniProjects/AppBlocker/cache/diegoibarra.AppBlocker.out</string>
	</dict>
</plist>