# OutlookCalendarUpdate/diegoibarra.OutlookCalendarUpdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>diegoibarra.OutlookCalendarUpdate</string>
		<key>Program</key>
		<string>/Users/diegoibarra/Developer/1_myProjects/MiniProjects/OutlookCalendarUpdate/CalendarUpdate</string>
		<key>WorkingDirectory</key>
		<string>/Users/diegoibarra/Developer/1_myProjects/MiniProjects/OutlookCalendarUpdate</string>
		<key>KeepAlive</key>
		<true/>
		<key>ThrottleInterval</key>
		<integer>1800</integer>
		<key>StandardErrorPath</key>
		<string>/Users/diegoibarra/Developer/1_myProjects/MiniProjects/OutlookCalendarUpdate/.cache/diegoibarra.OutlookCalendarUpdate.err</string>
		<key>StandardOutPath</key>
		<string>/Users/diegoibarra/Developer/1_myProjects/MiniProjects/OutlookCalendarUpdate/.cache/diegoibarra.OutlookCalendarUpdate.out</string>
	</dict>
</plist>